# deploy/systemd/library-prep-backup.timer
[Unit]
Description=Daily Library Prep PostgreSQL backup

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
RandomizedDelaySec=15m
Unit=library-prep-backup.service

[Install]
WantedBy=timers.target
